# jitsi-meet-ha-addon/supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/logs/supervisor.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ── Prosody (XMPP signaling) ─────────────────────────────────────────────────
[program:prosody]
command=/usr/bin/prosody
autostart=true
autorestart=true
priority=10
stdout_logfile=/data/logs/prosody.log
stderr_logfile=/data/logs/prosody.log
user=prosody

# ── Jicofo (conference focus) ────────────────────────────────────────────────
[program:jicofo]
command=/usr/share/jicofo/jicofo.sh
autostart=true
autorestart=true
priority=20
startsecs=5
startretries=10
stdout_logfile=/data/logs/jicofo.log
stderr_logfile=/data/logs/jicofo.log

# ── Jitsi Video Bridge (media relay) ────────────────────────────────────────
[program:jvb]
command=/usr/share/jitsi-videobridge/jvb.sh
autostart=true
autorestart=true
priority=20
startsecs=5
stdout_logfile=/data/logs/jvb.log
stderr_logfile=/data/logs/jvb.log

# ── Nginx (web frontend) ─────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=30
startsecs=3
stdout_logfile=/data/logs/nginx.log
stderr_logfile=/data/logs/nginx.log
